# supervisord.conf
[unix_http_server]
file = /tmp/supervisord/supervisor.sock

[supervisorctl]
serverurl = unix:///tmp/supervisord/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
logfile=%(ENV_HOME)s/logs/supervisord.log ; supervisord log file
logfile_maxbytes=50MB       ; maximum size of logfile before rotation
logfile_backups=10          ; number of backed up logfiles
loglevel=debug              ; info, debug, warn, trace
pidfile=/tmp/supervisord/supervisord.pid ; pidfile location
childlogdir=%(ENV_HOME)s/logs/            ; where child log files will live

[program:ssh-agent]
command=bash -c "rm /tmp/ssh-agent.sock -f && /usr/bin/ssh-agent -d -a /tmp/ssh-agent.sock"
priority=1
autorestart=true


[program:wsgi-app]
directory=/opt/marie-icr
command=python wsgi.py
process_name=%(program_name)s
user=root
autostart=true
autorestart=true
redirect_stderr=true
redirect_stderr=true


[program:registry]
directory=/opt/marie-icr
command=python ./register.py --config /etc/marie/marie.yml
process_name=%(program_name)s
user=app-svc
autostart=true
autorestart=true
redirect_stderr=true
redirect_stderr=true


